# systemd/.config/systemd/user/gpg-agent.service
[Unit]
Description=GPG Agent

[Service]
ExecStart=/usr/bin/gpg-agent --server

[Install]
WantedBy=default.target
