# systemd/.config/systemd/user/mpdris2.service
[Unit]
Description=MPD + Mpris interop

[Service]
ExecStart=/usr/local/bin/mpDris2

[Install]
WantedBy=default.target
