# systemd/.config/systemd/user/twmnd.service
[Unit]
Description=TWMN Daemon

[Service]
ExecStart=/home/matus/.local/bin/twmnd

[Install]
WantedBy=default.target
